# sys/etc/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.
/usr/bin/kernelTune

touch /var/lock/subsys/local
iwconfig wlan0 rate 54M
/usr/bin/ssdTune
#./mountCompressed.sh

